# home/.config/systemd/user/user@.service.d/environment.conf
[Service]
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%I/dbus/user_bus_socket
Environment=XDG_SEAT=seat0
Environment=XDG_VTNR=1

