# home/.config/systemd/user/bspwm.service
[Unit]
Description=bspwm
Requires=xorg.service
After=xorg.service

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/usr/bin/bspwm 

