# home/.config/systemd/user/bspwm.service
[Unit]
Requires=xorg@0.service
After=xorg@0.service

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/usr/bin/bspwm 

