# home/.config/systemd/user/environment-variables.service
[Service]
ExecStart=/usr/bin/xargs systemctl --user set-environment < %h/.envVariables
