# home/.config/systemd/user/environment.target
[Unit]
Requires=environment-variables.service
Requires=dbus.socket
